# templates/systemd/curator.timer.erb
[Unit]
Description=Elasticsearch Curator <%= @name %> Timer

[Timer]
OnCalendar=<%= @systemd_timer %>
Unit=curator_<%= @name %>.service

[Install]
WantedBy=multi-user.target
